# supervisor_phpfpm.conf
[program:phpfpm]
command=php5-fpm
autostart=true