# supervisor_nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true